# trouble.conf
description "Start the doubletrouble server"
author "The doubletrouble team"
setuid trouble

start on runlevel [2345]
exec pm2 start /home/trouble/doubletrouble/app.js

