# trouble.service
[Unit]
Description=Start the doubletrouble server

[Service]
Type=simple
ExecStart=/usr/bin/nodejs /home/trouble/doubletrouble/app.js

